# webcam.service
[Unit]
Description=Webcam Exhibit Fullscreen Display
After=multi-user.target
# Prevents getty from competing for VT1 (the display terminal)
Conflicts=getty@tty1.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
# Runs Python directly — no X server required.
# SDL_VIDEODRIVER=kmsdrm lets Pygame render to the framebuffer via
# the Linux KMS/DRM stack, which works without a display manager.
ExecStart=/usr/bin/python3 __INSTALL_DIR__/main.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=SDL_VIDEODRIVER=kmsdrm

[Install]
WantedBy=multi-user.target
